# conf/supervisor.conf
[program:<title>]
user=<user>
directory=/home/<user>/<folder>
command=/home/<user>/<folder>/venv/bin/gunicorn --worker-class=gevent --worker-connections=1000 --workers=1 --bind=127.0.0.1:<port> core.wsgi:application

autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/home/<user>/<folder>/logs/err.log
stdout_logfile=/home/<user>/<folder>/logs/out.log

[program:<title>_bot]
user=<user>
directory=/home/<user>/<folder>
command=/home/<user>/<folder>/venv/bin/python manage.py run_uvicorn

autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/home/<user>/<folder>/logs/bot_err.log
stdout_logfile=/home/<user>/<folder>/logs/bot_out.log


[program:<title>_jobs]
user=<user>
directory=/home/<user>/<folder>
command=/home/<user>/<folder>/venv/bin/python manage.py run_jobs

autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/home/<user>/<folder>/logs/jobs_err.log
stdout_logfile=/home/<user>/<folder>/logs/jobs_out.log


[program:<title>_celery]
user=<user>
directory=/home/<user>/<folder>
command=/home/<user>/<folder>/venv/bin/celery -A core worker --loglevel=info

autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
stderr_logfile=/home/<user>/<folder>/logs/celery_err.log
stdout_logfile=/home/<user>/<folder>/logs/celery_out.log
